# config/systemd/user/kidproxy.timer
[Unit]
Description=Update Kidproxy

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1min
AccuracySec=5s

[Install]
WantedBy=timers.target
