# config/systemd/user/kidproxy.timer
[Unit]
Description=Update Policy

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
